# .config/systemd/user/b3.service
[Unit]
Description=UrT 4.3 B3
After=urt43.service

[Install]
WantedBy=default.target

[Service]
Type=exec
WorkingDirectory=%h/b3
ExecStart=taskset -c 1 %h/b3/.venv/bin/python -B -u -m b3 -c %h/.config/b3/b3_ts.ini
TimeoutSec=600
Restart=on-failure
LimitNOFILE=5000
