# bin/urt43_scrim.service
# 30plus UrT systemd service file

[Unit]
Description=UrT 4.3 Scrim
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
User=urtadmin
Group=urtadmin
PIDFile=/game/run/urt43_scrim.pid
ExecStart=/game/bin/start_urt_scrim.sh
ExecStop=/game/bin/stop_urt_scrim.sh
TimeoutSec=30
Restart=on-failure
LimitNOFILE=5000
