# bin/urt43_b3_scrim.service
# 30plus B3 systemd service file

[Unit]
Description=UrT 4.3 B3 Scrim
After=network.target urt43_scrim.service

[Install]
WantedBy=multi-user.target

[Service]
User=urtadmin
Group=urtadmin
PIDFile=/game/run/b3_scrim.pid
ExecStart=/game/bin/start_b3_scrim.sh
ExecStop=/game/bin/stop_b3_scrim.sh
TimeoutSec=600
Restart=on-failure
LimitNOFILE=5000
